# back/action-scheduler.conf
[program:action-scheduler]
command=/path/to/rkllm-server/external/rag/back/venv/bin/python src/command/run_action_scheduler.py --interval 60
directory=/path/to/rkllm-server/external/rag/back
user=www-data
autostart=true
autorestart=true
redirect_stderr=true
stdout_logfile=/var/log/action-scheduler/action-scheduler.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=10
environment=PATH="/path/to/rkllm-server/external/rag/back/venv/bin"
